# install/ubuntu24/netalertx.service
[Unit]
Description=NetAlertX - Network, presence scanner and alert framework

[Service]
EnvironmentFile=/etc/default/netalertx
PassEnvironment=INSTALL_SYSTEM_NAME INSTALLER_DIR INSTALL_DIR PHPVERSION VIRTUAL_ENV PATH
ExecStartPre=/usr/bin/bash "${INSTALLER_DIR}/pre-start.sh"
ExecStart=/usr/bin/python3 "${INSTALL_DIR}/server"
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
