# install/alpine-docker/etc/crontabs/root
# Schedule cron jobs
* * * * * /app/back/cron_script.sh
