# install/crontab


# Schedule cron jobs
* * * * * /app/back/cron_script.sh
* * * * * echo "$(date +'%Y-%m-%d %H:%M:%S') - Cron job ran" >> /app/front/log/cron_timestamp.log
